# systemd/user/staff-briefs.service
[Unit]
Description=Office Runtime - Staff Briefs

[Service]
Type=oneshot
WorkingDirectory=/home/matias/repos/office-auto-lab
ExecStart=/usr/bin/bash /home/matias/repos/office-auto-lab/src/office_runtime/scripts/office_run.sh daily --scan-mode existing
